# linux/supervisord.virtualenv.ini
[program:nounou2gmaps]
user = photos
command=/home/photos/nounou2gmaps/env/bin/python /home/photos/nounou2gmaps/nounou2gmaps.py
directory=/home/photos/nounou2gmaps
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8,PATH="/home/photos/nounou2gmaps/env/bin/"
stdout_logfile=/var/log/photos/supervisor.nounou2gmaps.log
redirect_stderr = true
